# profile/airootfs/etc/systemd/system/docker.service.d/nfs-dependency.conf
[Unit]
After=mnt-synology-harbor_srv.mount
Requires=mnt-synology-harbor_srv.mount
